# dropin_service.service
[Service]
Environment=DOOMWADDIR=/usr/local/packages/xxxxxx